# packaging/systemd/tarisya-agent.service
[Unit]
Description=Tarisya Agent
Documentation=https://tarisya.nurfatkhur.com
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=tarisya
Group=tarisya
WorkingDirectory=/var/lib/tarisya-agent
EnvironmentFile=/etc/tarisya/agent.env
ExecStart=/usr/local/bin/tarisya-agent
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
UMask=0077

NoNewPrivileges=true
PrivateTmp=true
ProtectControlGroups=true
ProtectHome=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
